# rpm/doorcgi.conf
start on (net-device-up and local-filesystems and runlevel [2345])
stop on runlevel [016]
respawn

console none

env LANG=en_US.UTF-8

exec su -s /bin/bash -c '/usr/bin/doorcgi' nobody
